# plexmediaserver/files/plexmediaserver.in
#!/bin/sh
#
# PROVIDE: plexmediaserver
# KEYWORD: shutdown

. /etc/rc.subr

name=plexmediaserver
rcvar=plexmediaserver_enable
extra_commands="status"

start_cmd="plexmediaserver_start"
status_cmd="plexmediaserver_status"
stop_cmd="plexmediaserver_stop"

load_rc_config $name
: ${plexmediaserver_support_path="%%SUPPORT_PATH%%"}
: ${plexmediaserver_script_path="%%SCRIPT_PATH%%"}

PIDFILE=${plexmediaserver_support_path}/Plex\ Media\ Server/plexmediaserver.pid

plexmediaserver_start() {
	if [ -f "$PIDFILE" ]; then
		echo "Plex Media Server already running with pid: `cat "$PIDFILE"`"
	else
		echo "Starting Plex Media Server"
		${plexmediaserver_script_path}/start.sh ${plexmediaserver_support_path} > /dev/null 2>&1 &
	fi
}

plexmediaserver_status() {
	if [ -f "$PIDFILE" ]; then
		PID=`cat "$PIDFILE"`
		PGID=`ps -o "pgid=" $PID`
		if [ -z "$PGID" ]; then
			echo "Plex Media Server is not running"
		else
			echo "Plex Media Server is running with pid: `cat "$PIDFILE"`"
			ps -o "pid pgid command" | egrep "^[0-9]+\ $PGID\ "
		fi
	else
		echo "Plex Media Server is not running"
	fi
}

plexmediaserver_stop() {
	if [ -f "$PIDFILE" ]; then
		PID=`cat "$PIDFILE"`
		PID_CMD=`ps -o "command=" $PID`
		echo $PID_CMD | grep -q Plex
		if [ $? -eq 0 ]; then
			PGID=`ps -o "pgid=" $PID`
			echo "Stopping Plex Media Server"
			ps -o "pid pgid command" | egrep "^[0-9]+\ $PGID\ "
			kill -9 -- -$PGID
		else
			echo "Plex Media Server is not running. Pidfile removed"
		fi
		rm "$PIDFILE"
	else
		echo "Plex Medis Server is not running"
	fi
}

run_rc_command "$1"
