# files/plexmediaserver.in
#!/bin/sh
# Created by: Created by: KalleDK <plexmaintainer@k-moeller.dk>
#
# $FreeBSD$
#
# PROVIDE: plexmediaserver
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf to enable the Plex Media Server:
#
# plexmediaserver_enable="YES"
#

. /etc/rc.subr

name=plexmediaserver
rcvar=plexmediaserver_enable
load_rc_config $name

: ${plexmediaserver_enable:=NO}
: ${plexmediaserver_user:="%%USERS%%"}

extra_commands=status
start_precmd=plex_precmd
start_postcmd=plex_postcmd
start_cmd=plex_startcmd
status_cmd=plex_statuscmd
stop_cmd=plex_stopcmd
pidfile=/var/run/plex.pid

plex_precmd()
{
	export SUPPORT_PATH="%%DATADIR%%/plexdata"
	export HOME="${SUPPORT_PATH}/Plex Media Server"
	export SCRIPTPATH="%%DATADIR%%"
	export LD_LIBRARY_PATH="${SCRIPTPATH}"
	export PLEX_MEDIA_SERVER_HOME="${SCRIPTPATH}"
	export PLEX_MEDIA_SERVER_MAX_PLUGIN_PROCS=6
	export PLEX_MEDIA_SERVER_APPLICATION_SUPPORT_DIR=$SUPPORT_PATH
	export PATH="${SCRIPTPATH}/Resources/Python/bin:$PATH"
	export LC_ALL="C"
	export LANG="C"
	ulimit -s 3000
}

plex_startcmd()
{
	echo "Starting ${name}."
	/usr/sbin/daemon -u %%USERS%% -f '%%DATADIR%%/Plex Media Server'
}

plex_postcmd()
{
	if [ -f "%%DATADIR%%/plexdata/Plex Media Server/plexmediaserver.pid" ]; then
		ln -sf "%%DATADIR%%/plexdata/Plex Media Server/plexmediaserver.pid" /var/run/plex.pid
	fi
}

plex_statuscmd()
{
	if pgrep -F ${pidfile} >/dev/null; then
		echo "${name} is running with pid: $(cat ${pidfile})";
	else
		echo "${name} is not running.";
		exit 1;
	fi
}

plex_stopcmd()
{
	if pgrep -F "${pidfile}" >/dev/null; then
		mypid=$(cat ${pidfile});
		kill ${sig_stop} ${mypid};
		wait_for_pids ${mypid};
		pkill -9 -f plexmediaserver;
	else
		echo "${name} is not running.";
		exit 1;
	fi
}

run_rc_command "$1"
