# files/plexmediaserver.in
#!/bin/sh
# Created by: Created by: KalleDK <plexmaintainer@k-moeller.dk>
#
# $FreeBSD$
#
# PROVIDE: plexmediaserver
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf to enable the Plex Media Server:
#
# plexmediaserver_enable="YES"
#

. /etc/rc.subr

name=plexmediaserver
rcvar=plexmediaserver_enable
load_rc_config $name

: ${plexmediaserver_support_path="%%SUPPORT_PATH%%"}
: ${plexmediaserver_script_path="%%SCRIPT_PATH%%"}
: ${plexmediaserver_user="%%USERS%%"}

start_cmd=plexmediaserver_start
status_cmd=plexmediaserver_status
stop_cmd=plexmediaserver_stop
extra_commands=status
pidfile="${plexmediaserver_support_path}/Plex Media Server/plexmediaserver.pid"

plexmediaserver_start()
{
	if pgrep -F "${pidfile}" >/dev/null; then
		echo "${name} already running with pid: $(cat "$pidfile")"
	else
		echo "Starting ${name}"
		export HOME=$(eval printf ~${plexmediaserver_user})
		daemon -u ${plexmediaserver_user} -f ${plexmediaserver_script_path}/start.sh ${plexmediaserver_support_path}
	fi
}

plexmediaserver_status()
{
	if pgrep -F "${pidfile}" >/dev/null; then
		echo "${name} is running with pid: $(cat "${pidfile}")";
	else
		echo "${name} is not running.";
		exit 1;
	fi
}

plexmediaserver_stop()
{
        if pgrep -F "${pidfile}" >/dev/null; then
		mypid=$(cat "${pidfile}");
		kill ${sig_stop} ${mypid};
		wait_for_pids ${mypid};
	else
		echo "${name} is not running.";
		exit 1;
	fi
}

run_rc_command "$1"
